# ai_trading.service
[Unit]
Description=AI Trading System with Telegram Commands (Demo Account)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/quant_system_clean
ExecStart=/usr/bin/python3 ai_trading_system.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/quant_system_clean
Environment=OANDA_API_KEY=a3699a9d6b6d94d4e2c4c59748e73e2d-b6cbc64f16bcfb920e40f9117e66111a
Environment=OANDA_ENVIRONMENT=practice

[Install]
WantedBy=multi-user.target
