# automated_trading.service
[Unit]
Description=Automated Trading System (Demo Account)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/quant_system_clean
ExecStart=/usr/bin/python3 automated_trading_system.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/opt/quant_system_clean
Environment=OANDA_API_KEY=${OANDA_API_KEY}
Environment=OANDA_ENVIRONMENT=practice

[Install]
WantedBy=multi-user.target
