# automated_trading.service
[Unit]
Description=Automated Trading System (Demo Account)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspace
ExecStart=/usr/bin/python3 /workspace/automated_trading_system.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/workspace
# API credentials are managed in code/env; do not hardcode secrets here
Environment=OANDA_ENVIRONMENT=practice

[Install]
WantedBy=multi-user.target
